# mirai-agent/services/mirai-agent.service
[Unit]
Description=Mirai Autonomous AI Agent
Documentation=https://github.com/AgeeKey/Mirai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/mirai/mirai-agent
Environment="PATH=/root/mirai/mirai-agent/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/root/mirai/mirai-agent/.env

# Запуск
ExecStart=/root/mirai/mirai-agent/venv/bin/python3 /root/mirai/mirai-agent/main.py

# Автоматический перезапуск при сбоях
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Безопасность
NoNewPrivileges=true
PrivateTmp=true

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mirai-agent

# Ресурсы
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
